# scripts/systemd/cdd-app.service
[Unit]
Description=Laravel CDD Application Server
After=network.target

[Service]
User=vircom
Group=vircom
WorkingDirectory=/home/vircom/.gemini/antigravity/scratch/cdd_app
ExecStart=/usr/bin/php artisan serve --host=0.0.0.0 --port=8001
Restart=always
RestartSec=5
StandardOutput=append:/home/vircom/.gemini/antigravity/scratch/cdd_app/storage/logs/app-server.log
StandardError=append:/home/vircom/.gemini/antigravity/scratch/cdd_app/storage/logs/app-server.log

[Install]
WantedBy=multi-user.target
